# rootfs/etc/s6-overlay/s6-rc.d/ssh-agent/run
#!/bin/bash
set -eao pipefail
echo "ssh-agent: INFO: Starting ssh-agent service..."
ssh-agent -D > /etc/ssh/.ssh-agent
