# rootfs/etc/s6-overlay/s6-rc.d/ssh-agent/run
#!/bin/bash
set -eao pipefail
SSH_AUTH_SOCK="${SSH_AUTH_SOCK:-/tmp/ssh-agent.sock}"
echo "ssh-agent: INFO: Starting ssh-agent service (SSH_AUTH_SOCK=${SSH_AUTH_SOCK})..."
ssh-agent -D -a ${SSH_AUTH_SOCK} >/dev/null
